# IOT/Virtual/Sec_iot/systemdService/iot-controller.service
[Unit]
Description=Analisador/Controlador de trafego IOT
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/lib/iot-server/
ExecStart=/usr/lib/iot-server/main.py

[Install]
WantedBy=multi-user.target